# tux.service
[Unit]
Description=Server web to aniministrate cargo project
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/tux %h
[Install]
WantedBy=default.target
