# tux.service
[Unit]
Description=Server web to aniministrate cargo project
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/tux
ExecStartPost=/usr/bin/notify-send "tux are running"
[Install]
WantedBy=default.target
