# tux.service
[Unit]
Description=Server web to aniministrate cargo project
After=network.target
[Service]
Type=simple
ExecStart=/usr/bin/tux
[Install]
WantedBy=graphical-session.target